# skills/research-colony/com.marvin.research-colony.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marvin.research-colony</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/gileskayo/.agents/venv/bin/python</string>
        <string>/Users/gileskayo/.agents/skills/research-colony/scripts/run_colony.py</string>
    </array>

    <!-- Run every day at 09:00 (after daily-digest at 08:30) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gileskayo/.claude/logs/research-colony.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gileskayo/.claude/logs/research-colony-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/gileskayo</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/gileskayo/.agents/venv/bin</string>
    </dict>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
