# skills/improve/com.marvin.daily-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marvin.daily-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/gileskayo/.agents/venv/bin/python</string>
        <string>/Users/gileskayo/.agents/skills/improve/scripts/daily_digest.py</string>
    </array>

    <!-- Run every day at 08:30 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gileskayo/.claude/logs/daily-digest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gileskayo/.claude/logs/daily-digest-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/gileskayo</string>
        <key>PATH</key>
        <string>/Users/gileskayo/.local/bin:/usr/local/bin:/usr/bin:/bin:/Users/gileskayo/.agents/venv/bin</string>
    </dict>

    <!-- Restart if the job crashes -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
